# Driver/crontab
*/5 * * * * sudo sh /usr/WeatherServiceClient/autorun.sh