# RPi software/crontab
*/5 * * * * sudo /usr/Weather_service/autorun.sh