# RPi software/crontab
*/5 * * * * sudo sh /usr/Weather_service/autorun.sh
